# examples/service/resources/debian/postinst
#!/bin/sh
# postinst script for gosrvlibexample

set -e

update-rc.d -f gosrvlibexample defaults

#DEBHELPER#

exit 0
